# services/v4l2loopback.service
[Unit]
Description=v4l2loopback
After=network.target

[Service]
ExecStart=/usr/bin/sudo /usr/bin/modprobe v4l2loopback exclusive_caps=1 video_nr=40,50,51,60,61,70,71 card_label="OBS Virtual Camera","virtual_desk_vaughan_1","virtual_desk_vaughan_2","virtual_bed_tripod_1","virtual_bed_tripod_2","virtual_bed_overhead_1","virtual_bed_overhead_2"

[Install]
WantedBy=default.target
